# ansible/roles/tang_server/templates/tangd-update.path.j2
# {{ ansible_managed }}
[Unit]
DefaultDependencies=no
Requires=opt.automount
After=opt.automount

[Path]
PathChanged={{ tang_server_db_dir }}
MakeDirectory=true
DirectoryMode=0700
