# ansible/templates/opt.automount.j2
# {{ ansible_managed }}
[Unit]
# containerd puts stuff in the opt directory by default
Before=containerd.service

[Automount]
Where=/opt

[Install]
WantedBy=multi-user.target
